# files/etc/hotplug.d/firewall/99-miniupnpd-restart
#!/bin/sh
[ "${ACTION}" = "reload" ] || [ "${ACTION}" = "restart" ] && {
    /etc/init.d/miniupnpd restart
}
EOF