# files/etc/hotplug.d/iface/51-upnp
#!/bin/ash
sleep 10
if [ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ]; then
    /overlay/update-upnp.sh
fi